# maria-backup-maker/config/crontab
# в”Њ minute (0вЂ“59)
# в”‚ в”Њ hour (0вЂ“23)
# в”‚ в”‚ в”Њ day of month (1вЂ“31)
# в”‚ в”‚ в”‚ в”Њ month (1вЂ“12)
# в”‚ в”‚ в”‚ в”‚ в”Њ day of week (0вЂ“7, Sun=0 or 7)
# в”‚ в”‚ в”‚ в”‚ в”‚
# * * * * * command_to_execute
00 00 * * * backup.sh
15 13 * * * backup.sh
15 18 * * * backup.sh
